# templates/systemd/pai-lite-sync.service
# pai-lite sync service for Linux (systemd)
#
# This service runs 'pai-lite tasks sync'.
# It's triggered by pai-lite-sync.timer.
#
# Installation:
#   1. Copy to ~/.config/systemd/user/
#   2. Copy pai-lite-sync.timer too
#   3. Run: systemctl --user enable pai-lite-sync.timer
#   4. Run: systemctl --user start pai-lite-sync.timer
#
# Or use: pai-lite triggers install

[Unit]
Description=pai-lite task sync

[Service]
Type=oneshot
# Update this path to your pai-lite installation
ExecStart=/usr/local/bin/pai-lite tasks sync
