# templates/systemd/pai-lite-mayor.timer
# pai-lite Mayor keepalive timer for Linux (systemd)
#
# Triggers pai-lite-mayor.service:
# - On boot/login (Persistent=true catches up if missed)
# - Every 15 minutes thereafter
#
# Installation:
#   1. Copy pai-lite-mayor.service and pai-lite-mayor.timer to ~/.config/systemd/user/
#   2. Run: systemctl --user enable --now pai-lite-mayor.timer
#
# Or use: pai-lite triggers install

[Unit]
Description=pai-lite Mayor keepalive timer

[Timer]
# Run on boot/login
OnBootSec=10s
# Then every 15 minutes
OnUnitActiveSec=15min
# Catch up if timer was missed (e.g., laptop was asleep)
Persistent=true

[Install]
WantedBy=timers.target
