# templates/systemd/pai-lite-mayor.service
# pai-lite Mayor keepalive service for Linux (systemd)
#
# Ensures the Mayor tmux session is running:
# - Starts on login (default.target)
# - Checks every 15 minutes via companion timer
# - Idempotent: exits early if already running
#
# Installation:
#   1. Copy pai-lite-mayor.service and pai-lite-mayor.timer to ~/.config/systemd/user/
#   2. Edit the ExecStart path if needed
#   3. Run: systemctl --user enable --now pai-lite-mayor.timer
#
# Or use: pai-lite triggers install

[Unit]
Description=pai-lite Mayor session keepalive
After=network.target

[Service]
Type=oneshot
# Update this path to your pai-lite installation
ExecStart=/usr/local/bin/pai-lite mayor start

[Install]
WantedBy=default.target
