# templates/systemd/pai-lite-sync.timer
# pai-lite sync timer for Linux (systemd)
#
# This timer triggers pai-lite-sync.service periodically.
#
# Installation:
#   1. Copy to ~/.config/systemd/user/
#   2. Copy pai-lite-sync.service too
#   3. Run: systemctl --user enable pai-lite-sync.timer
#   4. Run: systemctl --user start pai-lite-sync.timer
#
# Or use: pai-lite triggers install

[Unit]
Description=pai-lite task sync timer

[Timer]
# Run 60 seconds after boot
OnBootSec=60

# Then run every hour
OnUnitActiveSec=1h

# Persist timer across reboots
Persistent=true

[Install]
WantedBy=timers.target
