# templates/launchd/ludics-mag.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  ludics Mag keepalive agent

  Ensures the Mag tmux session is running:
  - Starts on login (RunAtLoad)
  - Checks every 15 minutes (StartInterval)
  - Idempotent: exits early if already running

  Installation:
    1. Copy to ~/Library/LaunchAgents/
    2. Edit the path to ludics
    3. Run: launchctl load ~/Library/LaunchAgents/ludics-mag.plist

  Or use: ludics triggers install
-->
<dict>
  <key>Label</key>
  <string>com.ludics.mag</string>

  <key>ProgramArguments</key>
  <array>
    <string>/path/to/ludics</string>
    <string>mag</string>
    <string>start</string>
  </array>

  <!-- Start on login -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Check every 15 minutes (900 seconds) -->
  <!-- Idempotent: ludics mag start exits early if session exists -->
  <key>StartInterval</key>
  <integer>900</integer>

  <!-- Environment -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>

  <!-- Logging -->
  <key>StandardOutPath</key>
  <string>/tmp/ludics-mag.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/ludics-mag.err.log</string>
</dict>
</plist>
